# system/.crontab
0 0 * * * /home/gavaga/.bin/update_readings
